# etc/pi/services/ros_screen.service
[Unit]
Description=start roscore
After=remote-fs.target
After=syslog.target
[Service]
ExecStart=/usr/local/bin/ros_screen.sh
Restart=always
RestartSec=60
[Install]
WantedBy=multi-user.target
